# Interfaces/2-Services/tcp_server_service_with_DBUS_interface/dashcam_tcp_server.service
[Unit]
Description=Dashcam TCP Server service

[Service]
Type=simple
ExecStart=/usr/bin/dashcam_tcp_server

[Install]
WantedBy=multi-user.target
