# Interfaces/2-Services/tcp_client_service_with_DBUS_interface/dashcam_tcp_client.service
[Unit]
Description=Dashcam TCP Client service

[Service]
Type=simple
ExecStart=/usr/bin/dashcam_tcp_client

[Install]
WantedBy=multi-user.target
